# config/people-counter.service
[Unit]
Description=People Counter Edge Pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/usr/src/people-counter
ExecStart=/usr/bin/python3 -m src.main --config /etc/people-counter/config.yaml
Restart=always
RestartSec=10
WatchdogSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=people-counter

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/people-counter /var/log/people-counter /etc/people-counter
PrivateTmp=true

[Install]
WantedBy=multi-user.target
