# config/people-counter-reset.service
[Unit]
Description=People Counter daily reset (dedup counters + counting totals)

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/usr/src/people-counter
Environment="PYTHONPATH=/usr/src/people-counter"
ExecStart=/usr/bin/python3 -c "from src.wifi_ble.dedup import DedupEngine; d = DedupEngine('/var/lib/people-counter/dedup.db'); d.reset_daily(); print('Dedup reset done')"
